# docker/php/fetch-rates.cron
*/5 * * * * www-data php /var/www/html/app/bin/console app:fetch-rates >> /var/log/fetch-rates.log 2>&1
