# docker/php/fetch-rates.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/5 * * * * root cd /var/www/html/app && /usr/local/bin/php bin/console app:fetch-rates >> /var/log/fetch-rates.log 2>&1

